# recipes-core/volatile-binds/files/volatile-binds.service.in
[Unit]
Description=Bind mount volatile @where@
DefaultDependencies=false
Before=local-fs.target
RequiresMountsFor=@what@ @where@
ConditionPathIsReadWrite=@rwdir@
ConditionPathIsReadWrite=!@rodir@

[Service]
Type=oneshot
RemainAfterExit=Yes
StandardOutput=syslog
TimeoutSec=0
ExecStart=/sbin/mount-copybind @what@ @where@
ExecStop=/sbin/umount @where@

[Install]
WantedBy=local-fs.target
